# .docker/images/php/worker/supervisor/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
username=dummy
password=dummy

[supervisord]
nodaemon=true
logfile=__SUPERVISOR_LOG/supervisord.log
pidfile=/var/run/supervisord.pid
user = __APP_USER_NAME

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=dummy
password=dummy

[inet_http_server]
port=__SUPERVISOR_PORT
username=__SUPERVISOR_USER_NAME
password=__SUPERVISOR_PASSWORD

[include]
files = /etc/supervisor/conf.d/*.conf

[program:worker]
process_name=%(program_name)s_%(process_num)02d
command=__PHP_WORKER_COMMAND
autostart=true
autorestart=true
user=__APP_USER_NAME
numprocs=__PHP_WORKER_PROCESSES_NUMBER
stopwaitsecs=120
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
