# .docker/images/php-nginx/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
username=dummy
password=dummy

[supervisord]
nodaemon=true
logfile=__SUPERVISOR_LOG/supervisord.log
pidfile=/var/run/supervisord.pid
user = __APP_USER_NAME

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=dummy
password=dummy

[inet_http_server]
port=__SUPERVISOR_PORT
username=__SUPERVISOR_USER_NAME
password=__SUPERVISOR_PASSWORD

[include]
files = /etc/supervisor/conf.d/*.conf

[program:php-fpm]
command=php-fpm82 -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0
